# curator4-1059849716.conf
[program:curator4-1059849716]
command=python3 /home/curator4/curator4-1059849716/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1059849716/test.err.log
stdout_logfile=/home/curator4/curator4-1059849716/test.out.log